# private_dot_config/private_systemd/private_user/Cloud-BorgBackup.service.tmpl
[Unit]
Description=Cloud BorgBackup Service

[Service]
Type=simple
ExecStart={{- .chezmoi.homdeDir -}}/.local/bin/backup/Cloud-BorgBackup
WorkingDirectory={{- .chezmoi.homdeDir -}}/.local/bin/backup
