# private_dot_config/private_systemd/private_user/megasync-autostart.service
[Unit]
Description=Start megasync at login with a small 6 second delay
After=graphical-session.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 6
ExecStart=megasync

[Install]
WantedBy=graphical-session.target

